# backend/infrastructure/deployment_templates/chainlit-app.service
[Unit]
Description=AI Booking - Application
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/AI_Property_Booking_concierge/frontend
Environment="PATH=/root/AI_Property_Booking_concierge/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CHAINLIT_APP_ROOT=/root/AI_Property_Booking_concierge/frontend"
EnvironmentFile=/root/AI_Property_Booking_concierge/backend/.env
ExecStart=/root/AI_Property_Booking_concierge/backend/.venv/bin/chainlit run chainlit_app.py --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
