# backend/infrastructure/deployment_templates/chainlit-app.service
[Unit]
Description=AI Booking Concierge - Chainlit Application
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/hotel_booking/frontend
Environment="PATH=/root/hotel_booking/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CHAINLIT_APP_ROOT=/root/hotel_booking/frontend"
EnvironmentFile=/root/hotel_booking/.env
ExecStart=/root/hotel_booking/backend/.venv/bin/chainlit run chainlit_app.py --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
