# backend/infrastructure/deployment_templates/rust-gateway.service
[Unit]
Description=AI Booking Concierge - Rust Database Gateway
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/hotel_booking/backend/infrastructure/database/rust_db_gateway
Environment="PATH=/root/.cargo/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/hotel_booking/backend/infrastructure/database/rust_db_gateway/target/release/main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
